# servicios/arcom_update.timer
[Unit]
Description=Ejecutar arcom_update cada lunes

[Timer]
OnCalendar=Mon *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
